# php/php83/supervisord.conf
; ============================================================
; Supervisor 配置 - PHP 8.3
; 文档: http://supervisord.org/configuration.html
;
; 用法：
;   新增后台进程（队列、定时任务等），在本文件追加 [program:xxx] 段
;   重启容器后生效，无需重新构建镜像
; ============================================================

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true                          ; 前台运行（Docker 必须）
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ============================================================
; PHP-FPM 主进程（必须保留）
; ============================================================
[program:php-fpm]
command=php-fpm -F                     ; -F 前台运行
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/supervisor/php-fpm.stdout.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/php-fpm.stderr.log
stderr_logfile_maxbytes=5MB

; ============================================================
; 示例：Laravel Queue Worker（取消注释使用）
; ============================================================
; [program:laravel-worker]
; command=php /var/www/html/your-project/artisan queue:work --sleep=3 --tries=3 --timeout=90
; autostart=true
; autorestart=true
; numprocs=2                           ; 启动2个 worker 进程
; process_name=%(program_name)s_%(process_num)02d
; user=www-data
; redirect_stderr=true
; stdout_logfile=/var/log/supervisor/laravel-worker.log
; stdout_logfile_maxbytes=5MB

; ============================================================
; 示例：Laravel Scheduler（定时任务，取消注释使用）
; ============================================================
; [program:laravel-scheduler]
; command=php /var/www/html/your-project/artisan schedule:run --no-interaction
; autostart=true
; autorestart=true
; user=www-data
; redirect_stderr=true
; stdout_logfile=/var/log/supervisor/laravel-scheduler.log

; ============================================================
; 示例：Swoole HTTP Server（取消注释使用）
; ============================================================
; [program:swoole-server]
; command=php /var/www/html/your-project/server.php
; autostart=true
; autorestart=true
; user=www-data
; redirect_stderr=true
; stdout_logfile=/var/log/supervisor/swoole-server.log
